# packaging/macos/scripts/postinstall
#!/bin/bash
# macOS postinstall script for wda-agent.
#
# Creates the state/log directories, loads the launchd daemon, and
# kicks off the first run. Runs as root during installation.

set -e

LAUNCHD_LABEL="com.sn360.desktop-agent"
LAUNCHD_PLIST="/Library/LaunchDaemons/${LAUNCHD_LABEL}.plist"
SDA_HOME="/var/lib/sn360-desktop-agent"
SDA_LOG="/var/log/sn360-desktop-agent"
SDA_CONF="/etc/sn360-desktop-agent"

# Legacy wazuh-desktop-agent install paths + launchd label from pre-rename
# builds. Migrated on install so operators upgrading keep their state,
# and any leftover launchd job registered under the old label is unloaded.
LEGACY_HOME="/var/lib/wazuh-desktop-agent"
LEGACY_LOG="/var/log/wazuh-desktop-agent"
LEGACY_CONF="/etc/wazuh-desktop-agent"
LEGACY_LABEL="com.wazuh.desktop-agent"
LEGACY_PLIST="/Library/LaunchDaemons/${LEGACY_LABEL}.plist"

migrate_legacy_dir() {
    src="$1"
    dst="$2"
    if [ -d "$src" ] && [ ! -e "$dst" ]; then
        mkdir -p "$(dirname "$dst")"
        mv "$src" "$dst"
        echo "migrated $src -> $dst" >&2
    fi
}

migrate_legacy_dir "$LEGACY_CONF" "$SDA_CONF"
migrate_legacy_dir "$LEGACY_HOME" "$SDA_HOME"
migrate_legacy_dir "$LEGACY_LOG"  "$SDA_LOG"

# Unload + remove the legacy launchd job so only the new label is registered.
if launchctl list "$LEGACY_LABEL" >/dev/null 2>&1; then
    launchctl unload "$LEGACY_PLIST" >/dev/null 2>&1 || true
fi
if [ -f "$LEGACY_PLIST" ]; then
    rm -f "$LEGACY_PLIST"
fi

mkdir -p "$SDA_HOME" "$SDA_LOG" "$SDA_CONF" "$SDA_CONF/sca"
chown -R root:wheel "$SDA_HOME" "$SDA_LOG"
chmod 0750 "$SDA_HOME" "$SDA_LOG"
chmod 0755 "$SDA_CONF" "$SDA_CONF/sca"
if [ -f "$SDA_CONF/config.yaml" ]; then
    chmod 0644 "$SDA_CONF/config.yaml"
fi

# Unload any previous version before reloading, so upgrades pick up
# plist changes immediately.
if launchctl list "$LAUNCHD_LABEL" >/dev/null 2>&1; then
    launchctl unload "$LAUNCHD_PLIST" >/dev/null 2>&1 || true
fi
launchctl load -w "$LAUNCHD_PLIST"

exit 0
